# debian/zabbix-iostat/usr/lib/systemd/system/zabbix-iostat.service

[Unit]
Description=Zabbix iostat collector

[Service]
ExecStart=/usr/bin/zab-iostat-hires.pl

[Install]
WantedBy=multi-user.target

