# debian/zabbix-iostat/etc/cron.d/zabbix-iostat.cron

# Collect data iops for zabbix
* * * * *	root	/usr/bin/iostat-collect.sh /tmp/iostat.out 60

